# .config/systemd/user/wallpaper.timer
[Timer]
Unit=wallpaper.service
OnUnitActiveSec=1h

[Install]
WantedBy=default.target
WantedBy=timers.target

